# rc.local
#!/bin/bash
/usr/local/bin/install-shell
